# opt/etc/cron.daily/50-keenetic-pbr-lists-update.sh
#!/opt/bin/sh

logger -t "keenetic-pbr" "Updating lists"
/opt/usr/bin/keenetic-pbr download