# etc/cron.daily/50-keen-pbr-lists-update.sh
#!/opt/bin/sh

/opt/etc/init.d/S80keen-pbr download-lists
/opt/etc/init.d/S80keen-pbr apply-lists