# packages/openwrt/keen-pbr/files/etc/hotplug.d/iface/70-keen-pbr
#!/bin/sh

if [ -x /etc/init.d/keen-pbr ] && /etc/init.d/keen-pbr enabled; then
    logger -t "keen-pbr" "Sending reload signal to keen-pbr for $INTERFACE due to $ACTION of $INTERFACE${DEVICE:+ (${DEVICE})}"
    /etc/init.d/keen-pbr on_interface_reload "$INTERFACE" "$ACTION"
fi
