# keenetic-pbr/opt/etc/cron.daily/50-keenetic-pbr-lists-update.sh
#!/opt/bin/sh

/opt/etc/init.d/S80keenetic-pbr download-lists
/opt/etc/init.d/S80keenetic-pbr apply-lists