# packages/openwrt/keen-pbr/files/etc/hotplug.d/firewall/70-keen-pbr
#!/bin/sh

[ "$ACTION" = "reload" ] || [ "$ACTION" = "restart" ] || exit 0

if [ -x /etc/init.d/keen-pbr ] && /etc/init.d/keen-pbr enabled; then
    logger -t "keen-pbr" "Sending SIGUSR1 to keen-pbr due to firewall action: $ACTION"
    /etc/init.d/keen-pbr reapply_firewall
fi
