# packages/openwrt/keen-pbr/files/etc/hotplug.d/firewall/70-keen-pbr
#!/bin/sh

exec /usr/lib/keen-pbr/firewall.sh
